# .ansible/collections/ansible_collections/rridane/base_systems/roles/system_manage_cntlm/templates/cntlm.service.j2
[Unit]
Description=Cntlm NTLM/NTLMv2 authenticating HTTP proxy
After=network.target

[Service]
ExecStart={{ cntlm_bin_path }} -c {{ cntlm_conf_path }} -v
Restart=always

[Install]
WantedBy=multi-user.target
